# bingPyWallpaper.conf
start on runlevel [2345]
stop on runlevel [!2345]

exec /home/sumedh/Projects/BingPyWallpaper/bingWallPaper.py